# configs/systemd/random-omarchy-theme.service
[Unit]
Description=Apply random Omarchy theme (oneshot)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/random-omarchy-theme.sh

